# openwrt/etc/uci-defaults/commotiond
#!/bin/sh
[ -f "/usr/bin/patch" ] && { \
  cd /etc/init.d/
  patch -p0 < /usr/share/commotion/patches/dnsmasq.patch
  
  cd /lib/wifi/
  patch -p0 < /usr/share/commotion/patches/mac80211.sh.patch
}

uci add ucitrack postcommotion
uci set ucitrack.@postcommotion[0].init='network'
uci set ucitrack.@network[0].init='commotiond'
uci add_list ucitrack.@network[0].affects=postcommotion
uci commit ucitrack

exit 0
