# openwrt/etc/uci-defaults/commotiond
#!/bin/sh
[ -f "/usr/bin/patch" ] && { \
  cd /etc/init.d/
  patch -p0 < /usr/share/commotion/patches/dnsmasq.patch
}

exit 0
