# openwrt/etc/uci-defaults/commotiond
#!/bin/sh

. /lib/config/uci.sh

[ -f "/etc/crontabs/root" ] || touch /etc/crontabs/root

grep -q "commotiond" /etc/crontabs/root || { \
  echo '* * * * *   /etc/init.d/commotiond cron' >> /etc/crontabs/root
  /etc/init.d/cron restart
}

[ -z $(uci_get ucitrack @setup_wizard[0]) ] && { \
  uci_add ucitrack setup_wizard
  uci_set ucitrack @setup_wizard[0] init commotiond
  uci_commit ucitrack
}

exit 0
