# openwrt/etc/uci-defaults/commotiond
#!/bin/sh
[ -f "/usr/bin/patch" ] && { \
  cd /etc/init.d/
  patch -p0 < /usr/share/commotion/patches/dnsmasq.patch
  
  cd /lib/wifi/
  patch -p0 < /usr/share/commotion/patches/mac80211.sh.patch
}

echo '192.168.1.20  thisnode' >> /etc/hosts

exit 0
